# embedded/files/etc/init.d/neto
#!/bin/sh /etc/rc.common

START=92
STOP=12
USE_PROCD=1

CONFIG_NAME="neto"
NETOD="/usr/bin/netod"
DNSMASQ_SNIPPET="neto.conf"
DNSMASQ_STATE_DIR="/var/lib/neto"
DNSMASQ_NORESOLV_STATE="$DNSMASQ_STATE_DIR/dnsmasq-noresolv.prev"
DNSMASQ_ADDSUBNET_STATE="$DNSMASQ_STATE_DIR/dnsmasq-addsubnet.prev"
CRON_FILE="/etc/crontabs/root"
CRON_BEGIN="# neto subscriptions begin"
CRON_END="# neto subscriptions end"

neto_reload_dnsmasq() {
	[ -x /etc/init.d/dnsmasq ] || return 0
	/etc/init.d/dnsmasq reload >/dev/null 2>&1 || /etc/init.d/dnsmasq restart >/dev/null 2>&1 || true
}

neto_reload_cron() {
	[ -x /etc/init.d/cron ] || return 0
	/etc/init.d/cron reload >/dev/null 2>&1 || /etc/init.d/cron restart >/dev/null 2>&1 || true
}

neto_remove_cron() {
	local tmp
	[ -f "$CRON_FILE" ] || return 0
	tmp="/tmp/neto-cron.$$"
	awk -v begin="$CRON_BEGIN" -v end="$CRON_END" '
		$0 == begin { skip = 1; next }
		$0 == end { skip = 0; next }
		!skip { print }
	' "$CRON_FILE" >"$tmp"
	cat "$tmp" >"$CRON_FILE"
	rm -f "$tmp"
	neto_reload_cron
}

neto_append_subscription_cron() {
	local section="$1" enabled auto_update hour
	config_get_bool enabled "$section" enabled 1
	config_get_bool auto_update "$section" auto_update 0
	[ "$enabled" -eq 1 ] && [ "$auto_update" -eq 1 ] || return 0
	case "$section" in
		""|*[!A-Za-z0-9_]*)
			return 0
			;;
	esac
	config_get hour "$section" update_hour "0"
	case "$hour" in
		""|*[!0-9]*)
			hour="0"
			;;
	esac
	[ "$hour" -ge 0 ] && [ "$hour" -le 23 ] || hour="0"
	printf "0 %s * * * /usr/bin/netod subscriptions update %s >/dev/null 2>&1 && /etc/init.d/neto restart >/dev/null 2>&1\n" "$hour" "$section" >>"$NETO_CRON_LINES"
}

neto_append_provider_cron() {
	local section="$1" enabled auto_update hour
	config_get_bool enabled "$section" enabled 1
	config_get_bool auto_update "$section" auto_update 0
	[ "$enabled" -eq 1 ] && [ "$auto_update" -eq 1 ] || return 0
	case "$section" in
		""|*[!A-Za-z0-9_]*)
			return 0
			;;
	esac
	config_get hour "$section" update_hour "0"
	case "$hour" in
		""|*[!0-9]*)
			hour="0"
			;;
	esac
	[ "$hour" -ge 0 ] && [ "$hour" -le 23 ] || hour="0"
	printf "5 %s * * * /usr/bin/netod providers update %s >/dev/null 2>&1 && /etc/init.d/neto restart >/dev/null 2>&1\n" "$hour" "$section" >>"$NETO_CRON_LINES"
}

neto_write_cron() {
	local tmp lines
	mkdir -p "$(dirname "$CRON_FILE")"
	touch "$CRON_FILE"
	tmp="/tmp/neto-cron.$$"
	lines="/tmp/neto-cron-lines.$$"
	awk -v begin="$CRON_BEGIN" -v end="$CRON_END" '
		$0 == begin { skip = 1; next }
		$0 == end { skip = 0; next }
		!skip { print }
	' "$CRON_FILE" >"$tmp"
	: >"$lines"
	NETO_CRON_LINES="$lines"
	config_foreach neto_append_subscription_cron subscription
	config_foreach neto_append_provider_cron provider
	if [ -s "$lines" ]; then
		printf "%s\n" "$CRON_BEGIN" >>"$tmp"
		cat "$lines" >>"$tmp"
		printf "%s\n" "$CRON_END" >>"$tmp"
	fi
	cat "$tmp" >"$CRON_FILE"
	rm -f "$tmp" "$lines"
	neto_reload_cron
}

neto_write_dnsmasq() {
	local dns_listen manage host port server old_noresolv old_addsubnet
	config_get_bool manage main manage_dnsmasq 1
	[ "$manage" -eq 1 ] || return 0
	config_get dns_listen main dns_listen "127.0.0.1:5353"
	host="${dns_listen%:*}"
	port="${dns_listen##*:}"
	[ -n "$host" ] && [ -n "$port" ] || return 1
	server="$host#$port"

	mkdir -p "$DNSMASQ_STATE_DIR"
	if [ ! -f "$DNSMASQ_NORESOLV_STATE" ]; then
		old_noresolv="$(uci -q get dhcp.@dnsmasq[0].noresolv || true)"
		if [ -n "$old_noresolv" ]; then
			printf "%s\n" "$old_noresolv" >"$DNSMASQ_NORESOLV_STATE"
		else
			printf "__missing__\n" >"$DNSMASQ_NORESOLV_STATE"
		fi
	fi
	if [ ! -f "$DNSMASQ_ADDSUBNET_STATE" ]; then
		old_addsubnet="$(uci -q get dhcp.@dnsmasq[0].addsubnet || true)"
		if [ -n "$old_addsubnet" ]; then
			printf "%s\n" "$old_addsubnet" >"$DNSMASQ_ADDSUBNET_STATE"
		else
			printf "__missing__\n" >"$DNSMASQ_ADDSUBNET_STATE"
		fi
	fi

	uci -q del_list dhcp.@dnsmasq[0].server="$server" || true
	uci add_list dhcp.@dnsmasq[0].server="$server" >/dev/null
	uci set dhcp.@dnsmasq[0].noresolv='1'
	uci set dhcp.@dnsmasq[0].addsubnet='32'
	uci commit dhcp
	rm -f /tmp/dnsmasq.d/$DNSMASQ_SNIPPET /etc/dnsmasq.d/$DNSMASQ_SNIPPET
	neto_reload_dnsmasq
}

neto_remove_dnsmasq() {
	local dns_listen host port server old_noresolv old_addsubnet
	config_get dns_listen main dns_listen "127.0.0.1:5353"
	host="${dns_listen%:*}"
	port="${dns_listen##*:}"
	server="$host#$port"

	uci -q del_list dhcp.@dnsmasq[0].server="$server" || true
	if [ -f "$DNSMASQ_NORESOLV_STATE" ]; then
		old_noresolv="$(cat "$DNSMASQ_NORESOLV_STATE")"
		if [ "$old_noresolv" = "__missing__" ]; then
			uci -q delete dhcp.@dnsmasq[0].noresolv || true
		else
			uci set dhcp.@dnsmasq[0].noresolv="$old_noresolv"
		fi
		rm -f "$DNSMASQ_NORESOLV_STATE"
	fi
	if [ -f "$DNSMASQ_ADDSUBNET_STATE" ]; then
		old_addsubnet="$(cat "$DNSMASQ_ADDSUBNET_STATE")"
		if [ "$old_addsubnet" = "__missing__" ]; then
			uci -q delete dhcp.@dnsmasq[0].addsubnet || true
		else
			uci set dhcp.@dnsmasq[0].addsubnet="$old_addsubnet"
		fi
		rm -f "$DNSMASQ_ADDSUBNET_STATE"
	fi
	uci commit dhcp
	rm -f /tmp/dnsmasq.d/$DNSMASQ_SNIPPET /etc/dnsmasq.d/$DNSMASQ_SNIPPET
	neto_reload_dnsmasq
}

start_service() {
	config_load "$CONFIG_NAME"

	local enabled singbox_bin
	config_get_bool enabled main enabled 1
	if [ "$enabled" -ne 1 ]; then
		neto_remove_cron
		return 0
	fi
	config_get singbox_bin main singbox_bin "/usr/libexec/neto/sing-box"

	"$NETOD" check || return 1
	"$NETOD" compile || return 1
	"$NETOD" apply || return 1
	"$singbox_bin" check -c /tmp/neto/sing-box.json || return 1

	procd_open_instance netod
	procd_set_param command "$NETOD" run
	procd_set_param respawn 3600 5 5
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance

	procd_open_instance sing-box
	procd_set_param command "$singbox_bin" run -c /tmp/neto/sing-box.json
	procd_set_param respawn 3600 5 5
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance

	neto_write_dnsmasq
	neto_write_cron
}

stop_service() {
	config_load "$CONFIG_NAME"

	local mark table
	config_get mark main mark "0x101"
	config_get table main table "101"

	nft delete table inet neto 2>/dev/null || true
	while ip -4 rule del fwmark "$mark" table "$table" 2>/dev/null; do
		:
	done
	ip -4 route del local default dev lo table "$table" 2>/dev/null || true
	neto_remove_dnsmasq
	neto_remove_cron
}

reload_service() {
	stop_service
	start_service
}

service_triggers() {
	procd_add_reload_trigger "$CONFIG_NAME"
}
